# hotspot/hotspot-auto.service
[Unit]
Description=Auto hotspot checker

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hotspot-auto
